# timers/twitchPanels.timer
[Unit]
Description=Recreate panels on twitch.tv/jaidchen hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target