# timers/twitchPanels.timer
[Unit]
Description=Recreate panels on twitch.tv/jaidchen daily

[Timer]
OnCalendar=*-*-* 20:00:00

[Install]
WantedBy=timers.target